# util/docker/web/runit/rr/run
#!/bin/bash

cd /var/azuracast/rr
rr serve