# util/docker/web/runit/rr/run
#!/bin/bash

sudo -E -u azuracast rr /var/azuracast/rr/.rr.yaml